# deploy/tixaro-update.service.template
[Unit]
Description=Tixaro-Update aus dem Admin-Center
After=docker.service network-online.target

[Service]
Type=oneshot
User=__INSTALL_USER__
WorkingDirectory=__INSTALL_DIR__
ExecStart=/bin/sh __INSTALL_DIR__/deploy/update-host.sh __INSTALL_DIR__
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
